# systemd/dcvpn.service
[Unit]
Description=DCVPN FastAPI service (uvicorn)
After=network.target

[Service]
User=spraxxxvpn
Group=spraxxxvpn
WorkingDirectory=/opt/dcvpn
EnvironmentFile=/opt/dcvpn/.env
ExecStart=/opt/dcvpn/venv/bin/uvicorn app:app --host 127.0.0.1 --port 8000
Restart=on-failure
LimitNOFILE=4096
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/opt/dcvpn /var/www/wg/clients /etc/wireguard

[Install]
WantedBy=multi-user.target
